# config/kanata/com.user.kanata-tray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.kanata-tray</string>

    <key>ProgramArguments</key>
    <array>
        <string>__KANATA_TRAY_BIN__</string>
        <string>--log-level</string>
        <string>2</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>KANATA_TRAY_CONFIG_DIR</key>
        <string>__HOME__/.config/kanata/kanata-tray</string>
        <key>KANATA_TRAY_LOG_DIR</key>
        <string>__HOME__/Library/Logs/kanata-tray</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ThrottleInterval</key>
    <integer>20</integer>
</dict>
</plist>
